# factory/systemd/udev/me176c-factory-wifiaddr@.service.in
# SPDX-License-Identifier: MIT
# Copyright (C) 2019 lambdadroid

[Unit]
Description=Set WiFi MAC Address (%I)
Requires=me176c-factory.service
After=me176c-factory.service

[Service]
Type=oneshot
EnvironmentFile=@RUNTIME_DIRECTORY@/environment
ExecStart=@BINARY_DIRECTORY@/ip link set dev %i address $WIFIADDR
